# drivers/eub-utils/eub-i2c.service.in
[Unit]
Description = esrille unbrick i2c bridge service

[Service]
ExecStart=@@PREFIX@@/bin/eub_i2cattach
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
Type=simple
KillMode=process

[Install]
WantedBy=sysinit.target
